# resources/rc.local
#!/bin/bash

if ! grep noatime /mnt/temp/etc/fstab >/dev/null 2>&1;then
  echo "fixing fstab..."
  sed -i 's/\(.*\)ext4[[:space:]]\+errors=remount-ro[[:space:]]\+\(.*\)/\1ext4 noatime,nodiratime,errors=remount-ro \2/g' /etc/fstab
  sync
  mount -o remount /&
fi

#set local time for windows users
if hash timedatectl >/dev/null 2>&1;then
  timedatectl set-local-rtc 1
fi

#fix resolv.conf to point to connman's
/etc/init.d/connman >/dev/null 2>&1

echo "booting..."
